# scripts/launchd/me.coilysiren.coily-audit-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Mac launchd agent for the coily audit dashboard. The Linux side is
  systemd/coily-audit-dashboard.{service,timer} on kai-server, firing every
  5 minutes; this is the Mac analog.

  Regenerates ~/.coily/dashboard.html every 5 minutes by shelling out to
  `coily audit dashboard --since 7d`. That path matches the coily default
  so `coily audit open` (no args) finds it.

  Install / uninstall via:
    scripts/coily-audit-dashboard-install-mac.py
    scripts/coily-audit-dashboard-install-mac.py --uninstall

  The installer rewrites {{HOME}} into the absolute path before copying to
  ~/Library/LaunchAgents/. Don't edit the copy in ~/Library/LaunchAgents/
  by hand - the source of truth lives here.
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>me.coilysiren.coily-audit-dashboard</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/coily</string>
      <string>audit</string>
      <string>dashboard</string>
      <string>--since</string>
      <string>7d</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/coily-audit-dashboard.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/coily-audit-dashboard.log</string>

    <key>ProcessType</key>
    <string>Background</string>
  </dict>
</plist>
