# systemd/claude-remote-control.service
[Unit]
Description=Claude Code remote-control daemon (kai-server, outbound HTTPS only)
After=syslog.target network.target nss-lookup.target network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
Restart=on-failure
RestartSec=10
User=kai
# projects/coilysiren is the parent of git repos, not a repo itself, so
# `--spawn worktree` cannot apply. All three hosts (kai-server, WSL,
# Windows-native) now operate on the projects/coilysiren parent dir
# with --spawn same-dir for consistency.
WorkingDirectory=/home/kai/projects/coilysiren
# `claude` is the npm-installed CLI under nvm. Source nvm.sh to pick up
# whichever node version nvm has set as default for kai, then exec so
# systemd tracks the claude process directly (no bash wrapper holding
# MAINPID). `claude login` must have been run once as kai with the
# active claude.ai subscription.
Environment=HOME=/home/kai
Environment=NVM_DIR=/home/kai/.nvm
# `--name` labels the pre-created session inside the daemon.
# `--remote-control-session-name-prefix` drives the bottom line of the
# claude.ai/code Remote Control dropdown row (default is `hostname`,
# which happens to be `kai-server` here but we set it explicitly so the
# dropdown stays correct if the host is ever renamed).
ExecStart=/bin/bash -c 'source $NVM_DIR/nvm.sh && exec claude remote-control --spawn same-dir --name kai-server --remote-control-session-name-prefix kai-server'
ExecStop=/bin/kill -TERM $MAINPID

[Install]
WantedBy=multi-user.target
