# systemd/claude-remote-control.service
[Unit]
Description=Claude Code remote-control daemon (kai-server, outbound HTTPS only)
After=syslog.target network.target nss-lookup.target network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
Restart=on-failure
RestartSec=10
User=kai
WorkingDirectory=/home/kai
# `claude` is the npm-installed CLI under nvm. Source nvm.sh to pick up
# whichever node version nvm has set as default for kai, then exec so
# systemd tracks the claude process directly (no bash wrapper holding
# MAINPID). `claude login` must have been run once as kai with the
# active claude.ai subscription.
Environment=HOME=/home/kai
Environment=NVM_DIR=/home/kai/.nvm
ExecStart=/bin/bash -c 'source $NVM_DIR/nvm.sh && exec claude remote-control --spawn worktree --name kai-server'
ExecStop=/bin/kill -TERM $MAINPID

[Install]
WantedBy=multi-user.target
