# systemd/claude-remote-control-restart.service
[Unit]
Description=Daily restart for claude-remote-control.service
After=claude-remote-control.service coilysiren-pull-all.service

[Service]
Type=oneshot
# Gate the restart on coilysiren-pull-all succeeding (and not still
# running). On a failed pull, skip silently so the daemon keeps serving
# yesterday's tree rather than restarting into a stale-or-broken state.
# Recovery: investigate why pull-all is unhappy, then `systemctl start
# claude-remote-control-restart.service` once it's green.
# See coilyco-flight-deck/infrastructure#211, coilyco-bridge/agentic-os-kai#612.
ExecCondition=/home/kai/projects/coilysiren/infrastructure/scripts/claude-remote-control-restart-precheck.sh
ExecStart=/usr/bin/systemctl restart claude-remote-control.service

[Install]
WantedBy=multi-user.target
