# systemd/claude-remote-control-wsl.service
[Unit]
Description=Claude Code remote-control daemon (kai-desktop-tower WSL)
After=syslog.target network.target nss-lookup.target network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
Restart=on-failure
RestartSec=10
User=kai
# /mnt/x is a DrvFs mount of the Windows X: drive; the same physical
# tree is also driven by the Windows-native Claude under
# X:\projects-x\coilysiren. Two distinct `--name`s keep the
# remote-control dropdown unambiguous.
WorkingDirectory=/mnt/x/projects-x/coilysiren
Environment=HOME=/home/kai
# WSL host uses the Claude native installer (no nvm, no Node wrapper):
# stable shim at ~/.local/bin/claude points at the versioned binary
# under ~/.local/share/claude/versions/. Call the shim directly so
# systemd tracks the claude process without a bash wrapper.
Environment=PATH=/home/kai/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# `--spawn same-dir` because projects/coilysiren is the parent of git
# repos, not a repo itself; worktree spawn cannot apply here.
# `--remote-control-session-name-prefix` overrides the default
# (`hostname`), which inside WSL returns the Windows host's name and
# collides with the Windows-native daemon in the dropdown.
ExecStart=/home/kai/.local/bin/claude remote-control --spawn same-dir --name kai-desktop-tower-wsl --remote-control-session-name-prefix kai-desktop-tower-wsl
ExecStop=/bin/kill -TERM $MAINPID

[Install]
WantedBy=multi-user.target
