# systemd/coilysiren-forgejo-mirror.service
[Unit]
Description=Mirror recently-pushed coilysiren GitHub repos to Forgejo
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=kai
WorkingDirectory=/home/kai
ExecStart=/home/kai/projects/coilysiren/infrastructure/scripts/coilysiren-forgejo-mirror.sh

[Install]
WantedBy=multi-user.target
