# systemd/repo-recall-update.timer
[Unit]
Description=Daily tick for repo-recall-update.service
Requires=repo-recall-update.service

[Timer]
# 03:15 local time daily. Staggered behind coily-update.timer (03:00)
# so the two don't fight over brew's update lock.
OnCalendar=*-*-* 03:15:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
